# com.roo.trading-report-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roo.trading-report-push</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/andrewbuis/Claude/Projects/trading-info-page/push-report.sh</string>
  </array>

  <!-- Fires at 12:15 AM every night — 15 min after Claude writes the HTML at midnight -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/andrewbuis/Claude/Projects/trading-info-page/push.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/andrewbuis/Claude/Projects/trading-info-page/push.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
